# ins/installer/installer-1.3.0/bin/install/appservice.service
[Unit]
Description=sobeyhive Application host set Engine
Documentation=http://docs.docker.com
After=remote-fs.target network.target
Before=docker.service

[Service]
ExecStart=/etc/init.d/appservice.sh
ExecReload=/bin/kill -USR2 $MAINPID

[Install]
WantedBy=multi-user.target
